# lark-c360-ai-quota/launchd/com.claude.c360-ai-quota.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.claude.c360-ai-quota</string>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-c</string>
		<string>echo c2V0IC1lCmV4cG9ydCBQQVRIPSIvVXNlcnMveHFkbWFjbWluaW00Ly5ucG0tZ2xvYmFsL2JpbjovdXNyL2xvY2FsL2JpbjovdXNyL2JpbjovYmluOi91c3Ivc2Jpbjovc2JpbiIKZXhwb3J0IExBUktfQ0xJPSIvVXNlcnMveHFkbWFjbWluaW00Ly5ucG0tZ2xvYmFsL2Jpbi9sYXJrLWNsaSIKZXhwb3J0IEMzNjBfTk9USUZZX1VTRVJfSUQ9Im91XzhjNjMzODgxYzlmMmM4YTA5NDI4MTg1ZDQ1ZmE4MzRjIgpMT0c9Ii90bXAvYzM2MF9mdWxsXyQoL2Jpbi9kYXRlICslWSVtJWRfJUglTSVTKS5sb2ciCi91c3IvYmluL2NhZmZlaW5hdGUgLWkgL3Vzci9sb2NhbC9iaW4vbm9kZSAvVXNlcnMveHFkbWFjbWluaW00L0Rlc2t0b3AvZmVpc2h1X2NsYXVkZS8uYWdlbnRzL3NraWxscy9sYXJrLWMzNjAtYWktcXVvdGEvc2NyaXB0cy9jMzYwX2NvbGxlY3QubWpzIC0tY29uZmlnIC9Vc2Vycy94cWRtYWNtaW5pbTQvRGVza3RvcC9mZWlzaHVfY2xhdWRlLy5hZ2VudHMvc2tpbGxzL2xhcmstYzM2MC1haS1xdW90YS9jMzYwLmNvbmZpZy5qc29uIC0tbG9nLWZpbGUgIiRMT0ciIHx8IFJDPSQ/ClJDPSR7UkM6LTB9CgojIFN0YXR1cyBjbGFzc2lmaWNhdGlvbiAoUDAgZml4IDIwMjYtMDYtMzApOiBkaXN0aW5ndWlzaCBsb2ctbWlzc2luZyAvCiMgemVyby1jdXN0b21lcnMgLyBwYXJ0aWFsLWZhaWwgLyBjbGVhbi1kb25lIHNvIHRoZSBJTSBub3RpZnkgZG9lc24ndAojIGxvb2sgInN1Y2Nlc3NmdWwiIHdoZW4gdGhlIHNjcmlwdCBuZXZlciBwcm9kdWNlZCBhIGxvZy4Kc3RhdHVzPSJ1bmtub3duIgp0cz0iPyIKdHNfaHVtYW49Ij8iCnRvdGFsX2N1c3RvbWVycz0iPyIKb2s9MApmYWlsPTAKc2tpcD0wCmRvbmU9MApzdW1fdG90YWwzMGQ9MAptaWFvZGFfcHVzaGVkPTAKCmlmIFsgISAtZiAiJExPRyIgXTsgdGhlbgogIHN0YXR1cz0ibG9nX21pc3NpbmciCmVsc2UKICB0b3RhbF9jdXN0b21lcnM9JChncmVwIC1tMSAtb0UgImZvdW5kIFswLTldKyBjdXN0b21lcnMiICIkTE9HIiB8IGdyZXAgLW9FICJbMC05XSsiIHx8IGVjaG8gIj8iKQogIG9rPSQoZ3JlcCAtY0UgIl5cW1swLTldKy9bMC05XStcXSBPSyAiICIkTE9HIiB8fCBlY2hvIDApCiAgZmFpbD0kKGdyZXAgLWNFICJeXFtbMC05XSsvWzAtOV0rXF0gRkFJTCAiICIkTE9HIiB8fCBlY2hvIDApCiAgc2tpcD0kKGdyZXAgLWNFICJeXFtbMC05XSsvWzAtOV0rXF0gU0tJUCAiICIkTE9HIiB8fCBlY2hvIDApCiAgZG9uZT0kKGdyZXAgLW9FICJeXFtbMC05XSsvWzAtOV0rXF0iICIkTE9HIiB8IHRhaWwgLTEgfCBncmVwIC1vRSAiWzAtOV0rIiB8IGhlYWQgLTEpCiAgZG9uZT0ke2RvbmU6LTB9CiAgc3VtX3RvdGFsMzBkPSQoZ3JlcCAtRSAiXlxbWzAtOV0rL1swLTldK1xdIE9LICIgIiRMT0ciIHwgZ3JlcCAtb0UgInRvdGFsMzBkPVswLTldKyIgfCBzZWQgInMvLio9Ly8iIHwgYXdrICJ7cys9XCQxfSBFTkQge3ByaW50IHMrMH0iKQogIG1pYW9kYV9wdXNoZWQ9JChncmVwIC1vRSAibWlhb2RhPVswLTldKyIgIiRMT0ciIHwgc2VkICJzLy4qPS8vIiB8IGF3ayAie3MrPVwkMX0gRU5EIHtwcmludCBzKzB9IikKICB0cz0kKGJhc2VuYW1lICIkTE9HIiAubG9nIHwgc2VkIC1FICJzL2MzNjBfZnVsbF8vLyIpCiAgdHNfaHVtYW49JChkYXRlIC1qIC1mICIlWSVtJWRfJUglTSVTIiAiJHRzIiAiKyVZLSVtLSVkICVIOiVNIiAyPi9kZXYvbnVsbCB8fCBlY2hvICIkdHMiKQoKICBpZiBbICIkdG90YWxfY3VzdG9tZXJzIiA9ICIwIiBdIHx8IFsgIiR0b3RhbF9jdXN0b21lcnMiID0gIj8iIF07IHRoZW4KICAgIHN0YXR1cz0iZW1wdHkiCiAgZWxpZiBbICIkZG9uZSIgLWdlICIkdG90YWxfY3VzdG9tZXJzIiBdIDI+L2Rldi9udWxsICYmIFsgIiRmYWlsIiAtZXEgMCBdIDI+L2Rldi9udWxsOyB0aGVuCiAgICBzdGF0dXM9ImRvbmUiCiAgZWxpZiBbICIkZmFpbCIgLWd0IDAgXSAyPi9kZXYvbnVsbDsgdGhlbgogICAgc3RhdHVzPSJwYXJ0aWFsX2ZhaWwiCiAgZWxzZQogICAgc3RhdHVzPSJwYXJ0aWFsIgogIGZpCmZpCgpjYXNlICIkc3RhdHVzIiBpbgogIGRvbmUpICAgICAgICAgcHJlZml4PSLinIUg5a6M5oiQIjs7CiAgcGFydGlhbF9mYWlsKSBwcmVmaXg9IuKaoO+4jyDpg6jliIblpLHotKUiOzsKICBwYXJ0aWFsKSAgICAgIHByZWZpeD0i4pqg77iPIOmDqOWIhuWujOaIkCI7OwogIGVtcHR5KSAgICAgICAgcHJlZml4PSLimqDvuI8gMCDlrqLmiLciOzsKICBsb2dfbWlzc2luZykgIHByZWZpeD0i4p2MIOS4u+aKk+WPliBsb2cg57y65aSxIjs7CiAgKikgICAgICAgICAgICBwcmVmaXg9IuKdkyDmnKrnn6XnirbmgIEiOzsKZXNhYwoKQk9EWT0iJHByZWZpeCAqKkMzNjAgJHRzX2h1bWFuKiogwrcgdmlldz0kdG90YWxfY3VzdG9tZXJzIGRvbmU9JGRvbmUgb2s9JG9rIGZhaWw9JGZhaWwgc2tpcD0kc2tpcCBtaWFvZGE9JG1pYW9kYV9wdXNoZWQgcmM9JFJDIGxvZz0kTE9HIgoiJExBUktfQ0xJIiBpbSArbWVzc2FnZXMtc2VuZCAtLXVzZXItaWQgIiRDMzYwX05PVElGWV9VU0VSX0lEIiAtLW1hcmtkb3duICIkQk9EWSIgfHwgZWNobyAibm90aWZ5IGZhaWxlZCIKZXhpdCAkUkMK | /usr/bin/base64 -d | /bin/bash</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/tmp/c360-ai-quota-launchd.err</string>
	<key>StandardOutPath</key>
	<string>/tmp/c360-ai-quota-launchd.out</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>0</integer>
		<key>Minute</key>
		<integer>30</integer>
	</dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>C360_NOTIFY_USER_ID</key>
		<string>ou_8c633881c9f2c8a09428185d45fa834c</string>
		<key>HOME</key>
		<string>/Users/xqdmacminim4</string>
		<key>LARK_CLI</key>
		<string>/Users/xqdmacminim4/.npm-global/bin/lark-cli</string>
		<key>PATH</key>
		<string>/Users/xqdmacminim4/.npm-global/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
</dict>
</plist>
